# root/etc/s6-overlay/s6-rc.d/init-mod-code-server-python3-add-package/run
#!/usr/bin/with-contenv bash

echo "*** Installing UV ***"
curl -LsSf https://astral.sh/uv/install.sh | sh
echo "*** UV installed ***"